# rakitanmanager/init.d/rakitanmanager
#!/bin/sh /etc/rc.common
# (C) 2025 rtaserver
# RakitanManager service control script

START=99
STOP=15
USE_PROCD=1

NAME="rakitanmanager"
DAEMON="/usr/share/rakitanmanager/core-manager.sh"

boot() {
    config_load rakitanmanager
    local enabled
	config_get_bool enabled "cfg" "enabled" 0

    if [ "$enabled" = 1 ]; then
        echo "Start Rakitan Manager"
        procd_open_instance rakitanmanager

        procd_set_param command "$DAEMON" -s
        procd_set_param pidfile "/var/run/rakitanmanager.pid"
        procd_set_param respawn
        procd_close_instance
    fi
}

start_service() {
    config_load rakitanmanager
    local enabled
	config_get_bool enabled "cfg" "enabled" 0

    [ "$enabled" = "1" ] || {
        echo "Service is disabled in configuration"
        return 0
    }
    
    if pidof core-manager.sh >/dev/null; then
        echo "RakitanManager is already running"
        return 0
    fi
    
    echo "Start Rakitan Manager"
    procd_open_instance rakitanmanager

    procd_set_param command "$DAEMON" -s
    procd_set_param pidfile "/var/run/rakitanmanager.pid"
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    echo "Stopping RakitanManager..."
    
    # First try to stop gracefully
    if [ -f /var/run/rakitanmanager.pid ]; then
        local pid=$(cat /var/run/rakitanmanager.pid 2>/dev/null)
        if [ -n "$pid" ] && kill -0 "$pid" 2>/dev/null; then
            kill -TERM "$pid"
            sleep 2
        fi
    fi
    
    # Force kill if still running
    if pidof core-manager.sh >/dev/null; then
            $DAEMON -k 2>/dev/null || {
            killall -9 core-manager.sh 2>/dev/null
        }
        sleep 1
    fi
    
    # Clean up pid file
    rm -f /var/run/rakitanmanager.pid 2>/dev/null
}

restart() {
    stop_service
    sleep 2
    start_service
}